# usr/lib/udev/rules.d/100-xps-power-manager.rules
ACTION=="change", SUBSYSTEM=="power_supply", RUN+="/usr/share/xps/udev/hooks/power-manager"
